# Themes/normal/etc/rc.local
cp /etc/pioneos/os-release /etc/os-release