# Themes/normal/etc/rc.local
#!/bin/sh
cp /usr/share/pioneos/os-release /etc/os-release